# systemd/ap_mode.service
[Unit]
Description=Start Access Point Mode on Boot
After=network.target
Wants=network.target

[Service]
Type=oneshot
# Call the Python helper (explicit interpreter for reliability)
ExecStart=/usr/bin/python3 /usr/local/bin/ap_mode.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target